# softlinks/Library/LaunchAgents/heymatthew.sync-git.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    test with
    plutil ~/Library/LaunchAgents/*.plist
  -->
  <key>Label</key>
  <string>heymatthew.sync-git.plist</string>

  <key>ProgramArguments</key>
  <array>
    <string>$HOME/bin/sync-git</string>
  </array>

  <!-- hourly -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
