# softlinks/Library/LaunchAgents/heymatthew.myrepos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>nz.matthew.plist</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mr</string>
    <string>--jobs</string>
    <string>2</string>
    <string>fetch</string>
  </array>

  <!-- hourly -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/myrepos.log</string>

  <key>StandardOutPath</key>
  <string>/tmp/myrepos.log</string>
</dict>
</plist>
