# softlinks/Library/LaunchAgents/heymatthew.logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    test with
    plutil ~/Library/LaunchAgents/*.plist
  -->
  <key>Label</key>
  <string>heymatthew.logrotate.plist</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/sbin/logrotate</string>
    <string>-s</string>
    <string>$HOME/logrotate.status</string>
    <string>$HOME/dotfiles/config/logrotate.conf</string>
  </array>

  <!-- hourly -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
